# ops/launchd/com.nico.wizard-status-ping.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nico.wizard-status-ping</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/lifeos.nico/Nico/Scripts/wizard-status-ping.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/lifeos.nico</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/lifeos.nico/Nico/Logs/wizard-status-ping.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lifeos.nico/Nico/Logs/wizard-status-ping-err.log</string>
</dict>
</plist>
